# systemd/ucas-login.timer
[Unit]
Description=Timer for UCAS SRUN auto login.

[Timer]
OnCalendar=*:*:00
Persistent=true

[Install]
WantedBy=timers.target