# utils/openwrt_cron/crontabs
* * * * * /usr/bin/python /root/UCAS-srun-login-script/run.py > /dev/null 2>&1